# overlays/browndeer-coprthr-1.6.0/usr/local/browndeer/scripts/clprocd.in
#!@BASH@

# chkconfig: 2345 85 15

prefix=@prefix@
exec_prefix=@exec_prefix@

start() {
	@bindir@/clprocd > /dev/null &
}

restart() {
	killall clprocd 
	@bindir@/clprocd > /dev/null &
}

stop() {
	killall clprocd 
}

case "$1" in 
	start)
		start
		;;
	restart)
		restart
		;;
	stop)
		stop
		;;
	*)
		echo "usage: clprocd {start|restart|stop}"
		exit 1
esac

exit 0

